# crontab
0 9 * * MON /path/to/venv/bin/python backend/tasks/weekly_email.py
## renovacion de suscripciones anuales